# .config/systemd/user/battery-monitor.service
[Unit]
Description=Battery Monitor Service
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-monitor.sh
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
