# .config/systemd/user/battery-monitor.service
[Unit]
Description=Battery Monitor Service
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-monitor.sh
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
